# attic/stuff/rc.local.in
export SX1276_HIGH=1
( while true ; do broker_mqtts /etc/sentrifarm-rsmb-config.txt >> /tmp/broker.log 2>&1 ; sleep 2 ; done ) &  
( while true ; do sleep 6 ; sx1276_mqttsn_bridge /dev/spidev0.1 connect 1883 >> /tmp/relay.log 2>&1 ; done )&
exit 0
